# HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/rollback-20260625-143102-debian-12-20260625.manifest
002|remove_file|/etc/modprobe.d/hitss-hardening-squashfs.conf||
003|remove_file|/etc/modprobe.d/hitss-hardening-udf.conf||
004|remove_file|/etc/modprobe.d/hitss-hardening-hfs.conf||
005|remove_file|/etc/modprobe.d/hitss-hardening-hfsplus.conf||
006|remove_file|/etc/modprobe.d/hitss-hardening-jffs2.conf||
007|remove_file|/etc/modprobe.d/hitss-hardening-freevxfs.conf||
009|remove_file|/etc/modprobe.d/hitss-hardening-usb_storage.conf||
010|remove_file|/etc/modprobe.d/hitss-hardening-dccp.conf||
012|remove_file|/etc/modprobe.d/hitss-hardening-rds.conf||
013|remove_file|/etc/modprobe.d/hitss-hardening-tipc.conf||
075|restore_file|/etc/postfix/main.cf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_075_etc_postfix_main.cf.bak|
090|remove_file|/etc/cron.allow||
090|remove_file|/etc/cron.deny||
090|remove_file|/etc/at.allow||
090|remove_file|/etc/at.deny||
092|remove_file|/etc/sysctl.d/99-hitss-hardening.conf||
092|sysctl|net.ipv4.conf.all.send_redirects|1|
092|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_092_etc_sysctl.d_99-hitss-hardening.conf.bak|
092|sysctl|net.ipv4.conf.default.send_redirects|1|
094|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_094_etc_sysctl.d_99-hitss-hardening.conf.bak|
094|sysctl|net.ipv4.conf.all.accept_redirects|1|
094|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_094_etc_sysctl.d_99-hitss-hardening.conf.bak|
094|sysctl|net.ipv4.conf.default.accept_redirects|1|
094|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_094_etc_sysctl.d_99-hitss-hardening.conf.bak|
094|sysctl|net.ipv4.conf.all.secure_redirects|1|
094|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_094_etc_sysctl.d_99-hitss-hardening.conf.bak|
094|sysctl|net.ipv4.conf.default.secure_redirects|1|
095|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_095_etc_sysctl.d_99-hitss-hardening.conf.bak|
095|sysctl|net.ipv4.conf.all.rp_filter|0|
095|restore_file|/etc/sysctl.d/99-hitss-hardening.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_095_etc_sysctl.d_99-hitss-hardening.conf.bak|
095|sysctl|net.ipv4.conf.default.rp_filter|0|
109|file_meta|/etc/ssh|root:root:755|
109|file_meta|/etc/ssh/sshd_config|root:root:644|
111|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_111_etc_ssh_sshd_config.bak|
112|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_112_etc_ssh_sshd_config.bak|
113|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_113_etc_ssh_sshd_config.bak|
120|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_120_etc_ssh_sshd_config.bak|
120|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_120_etc_ssh_sshd_config.bak|
123|restore_file|/etc/ssh/sshd_config|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_123_etc_ssh_sshd_config.bak|
135|remove_file|/etc/security/pwquality.conf||
136|restore_file|/etc/security/pwquality.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_136_etc_security_pwquality.conf.bak|
137|restore_file|/etc/security/pwquality.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_137_etc_security_pwquality.conf.bak|
137|restore_file|/etc/security/pwquality.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_137_etc_security_pwquality.conf.bak|
141|restore_file|/etc/security/faillock.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_141_etc_security_faillock.conf.bak|
142|restore_file|/etc/security/faillock.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_142_etc_security_faillock.conf.bak|
143|restore_file|/etc/login.defs|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_143_etc_login.defs.bak|
144|restore_file|/etc/login.defs|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_144_etc_login.defs.bak|
151|restore_file|/etc/shells|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_151_etc_shells.bak|
153|remove_file|/etc/profile.d/99-hitss-umask.sh||
154|remove_file|/etc/profile.d/99-hitss-tmout.sh||
157|remove_file|/etc/sudoers.d/99-hitss-157||
159|remove_file|/etc/sudoers.d/99-hitss-159||
160|restore_file|/etc/pam.d/su|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_160_etc_pam.d_su.bak|
162|restore_file|/etc/systemd/journald.conf|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_162_etc_systemd_journald.conf.bak|
183|remove_file|/etc/cron.daily/hitss-aide-check||
187|file_meta|/etc/shadow|root:shadow:640|
188|file_meta|/etc/shadow-|root:shadow:640|
189|file_meta|/etc/gshadow-|root:shadow:640|
190|file_meta|/etc/gshadow|root:shadow:640|
193|restore_file|/etc/shells|/home/vagrant/UBUNTU20-24_DEBIAN11-13/HARDENING_SHELLSCRIPT_UBT20-24_DEB11-13/logs/remed/rollback/20260625-143102-debian-12-20260625_193_etc_shells.bak|
201|file_meta|/home/vagrant|vagrant:vagrant:700|
201|file_meta|/home/renato|renato:renato:755|
